# oauth-mail-proxy/rootfs/etc/services.d/oauth-mail-proxy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Ensure the share directory exists
mkdir -p /share

#insert the account options into the config file
CONFIG_FILE="/etc/oauth-mail-proxy/emailproxy.config"

# Check for integration-provided account config
if [ -f "/share/oauth-mail-accounts.ini" ]; then
    bashio::log.info "Appending integration-provided accounts to emailproxy.config"
    echo "" >> "$CONFIG_FILE"
    echo "" >> "$CONFIG_FILE"
    echo "[accounts]" >> "$CONFIG_FILE"
    cat "/share/oauth-mail-accounts.ini" >> "$CONFIG_FILE"
else
    bashio::log.warning "No accounts configured in addon or integration. Please configure accounts in the Home Assistant integration."
    bashio::exit.ok
fi

## Run your program
exec /usr/bin/oauth-mail-proxy
